# Software/deployment/PowerMonitor.service
[Unit]
Description=Power Monitor - pull power monitoring data from RPi daughter board, process and publish via MQTT
After=network.target

[Service]
WorkingDirectory=/var/opt/PowerMonitor
SyslogIdentifier=PowerMonitor
ExecStart=/bin/sh -c "exec sudo java -jar /opt/PowerMonitor/Software/deployment/PowerMonitor.jar"
User=root
Type=simple

[Install]
WantedBy=multi-user.target